# packaging/DEBIAN/postinst
#!/bin/sh
set -e

build_sunwait() {
  tmpdir=$(mktemp -d)
  trap 'rm -rf "$tmpdir"' EXIT

  apt-get install -y build-essential ca-certificates curl
  curl -fsSL "https://github.com/risacher/sunwait/archive/refs/heads/master.tar.gz" -o "$tmpdir/sunwait.tar.gz"
  tar -xzf "$tmpdir/sunwait.tar.gz" -C "$tmpdir"
  make -C "$tmpdir/sunwait-master"
  install -m 0755 "$tmpdir/sunwait-master/sunwait" /usr/local/bin/sunwait
}

if ! command -v sunwait >/dev/null 2>&1; then
  if command -v apt-get >/dev/null 2>&1; then
    if ! apt-get install -y sunwait; then
      build_sunwait
    fi
  else
    echo "sunwait is required but apt-get is unavailable." >&2
    exit 1
  fi
fi

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload || true
  systemctl enable --now monitor-sun.service || true
fi

exit 0
