# systemd/monitor-sun.service
[Unit]
Description=Dynamic monitor day/night via DDC/CI and sunwait
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/monitor-sun-loop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
